# rc.local
clear
toilet --gay -f small "PiSP os"
echo "version 0.3.6"
sleep 3
